# infra/mac/io.mirotic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd plist for mirotic worker (Mac native). Install: ./mac/install-launchd.sh -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.mirotic</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bun</string>
    <string>run</string>
    <string>src/index.ts</string>
    <string>worker</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- PATH thừa các bin Mac cần: gh, claude, git, docker, ssh, rsync -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/mirotic/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/mirotic/stderr.log</string>

  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
